# core/caddy/caddy.initd
#!/sbin/openrc-run

name=caddy
command="/usr/sbin/$name"
command_args="-agree -email $caddy_email -conf $caddyfile"
command_background="yes"

start_stop_daemon_args="--user caddy:caddy"
pidfile="/run/$name.pid"

depend() {
  use logger dns
  need net
  after firewall
}
