# core/runit-openrc/runit-openrc.initd
#!/sbin/openrc-run

name=runit
command="/sbin/runsvdir"
command_args="/etc/service"
command_background="yes"

pidfile="/run/$name.pid"

depend() {
	need net
	after firewall
}
